# aspects/systemd/files/.config/systemd/user/colorscheme.timer
[Unit]
Description=Run script periodically

[Timer]
OnCalendar=*:0/1
Persistent=true

[Install]
WantedBy=timers.target
