# aspects/systemd/files/.config/systemd/user/colorscheme.service
[Unit]
Description=Run my script

[Service]
ExecStart=/bin/bash -c 'source $HOME/.profile && colorscheme.lua'

[Install]
WantedBy=timers.target

