# aspects/systemd/files/.config/systemd/user/txtatelier.service
[Unit]
Description=Txtatelier File Sync Service
Documentation=https://github.com/namjul/txtatelier
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# Environment configuration
EnvironmentFile=-%h/.config/txtatelier/txtatelier.env

# Start the service
ExecStart=%h/.local/bin/mise exec -- txtatelier

# Restart behavior
Restart=on-failure
RestartSec=10s
TimeoutStopSec=30s

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=txtatelier

[Install]
WantedBy=default.target
