# aspects/dotfiles/files/.config/autostart/gromit-mpx.desktop
[Desktop Entry]
Type=Application
Exec=gromit-mpx
