# aspects/dotfiles/files/.config/autostart/gromit-mpx.desktop
[Desktop Entry]
Type=Application
Name=Gromit Mpx
Exec=gromit-mpx
